# roles/thepoint_website/templates/thepoint.socket
[Unit]
Description={{ description }} socket

[Socket]
ListenStream=/run/{{ user }}/socket

[Install]
WantedBy=sockets.target
