# roles/thepoint_website/files/crontab
0 22 * * * ${HOME}/venv/bin/thepoint-backup
0 14 * * * ${HOME}/venv/bin/thepoint-sendrosteremails
